# desktop/build/sshagent/zpass-agent.service
# ZPass SSH Agent —— systemd user unit
#
# 安装：
#   1. cp zpass-agent.service ~/.config/systemd/user/
#   2. cp zpass-agent.socket  ~/.config/systemd/user/
#   3. 把 /usr/local/bin/zpass-agent 改成实际安装路径
#   4. systemctl --user daemon-reload
#   5. systemctl --user enable --now zpass-agent.socket
#
# 与 .socket 文件配合使用即可获得 socket activation：
# zpass-agent 进程仅在第一个 SSH 客户端 connect 时启动，
# 空闲时进程不存在 → 0 MB 内存占用。
#
# 详见同目录 zpass-agent.socket。

[Unit]
Description=ZPass SSH agent (zero-knowledge SSH key broker)
Documentation=https://github.com/zerx-lab/zpass
After=graphical-session.target
PartOf=graphical-session.target

[Service]
# Type=simple：zpass-agent 在 stderr 输出 log 并阻塞在 signal.Notify 等待 SIGTERM
Type=simple

# 改成实际安装路径（与 GUI 同目录或 $PATH 中可执行）
ExecStart=/usr/local/bin/zpass-agent

# 子进程内部已经处理 SIGTERM 触发优雅退出
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=10

# 仅当 socket 有连接时启动；socket 关闭自动停止
Restart=on-failure
RestartSec=2s

# 资源限制 —— agent 是只读协议适配器，绝不该用很多内存
MemoryMax=64M
TasksMax=32

# 安全 sandbox（systemd >= 232）
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%t/zpass %h/.zpass
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=default.target
